# services/minecraft-admin.service
[Unit]
Description=Minecraft admin
After=network.target

[Service]
Type=simple
User=minecraft
WorkingDirectory=/opt/minecraft
EnvironmentFile=/opt/minecraft/.env
Environment="PATH=/home/minecraft/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/minecraft/tools/start_minecraft_admin.sh
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
